# ansible/roles/raid0/templates/raid0.service
[Unit]
Description=RAID0 array setup
DefaultDependencies=no
After=local-fs-pre.target
Before=local-fs.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash /usr/sbin/raid0.sh

[Install]
WantedBy=local-fs.target
