# ansible/roles/raid0/templates/raid0.service
[Unit]
Description=raid0

[Service]
Type=oneshot
ExecStart=/bin/bash /usr/sbin/raid0.sh

[Install]
WantedBy=multi-user.target
