# ansible/roles/raid0/templates/raid0.mount
[Unit]
Description=raid0
After=raid0.service

[Mount]
What=/dev/md/ephemeral
Where={{ raid0_mount_path }}
Type=xfs
Options=defaults

[Install]
WantedBy=multi-user.target
